# conf/production/supervisord.conf
; Webfaction supervisor config file.
 
[inet_http_server]                         ; inet (TCP) server setings
port=127.0.0.1:25863                    ; (ip_address:port specifier, *:port for all iface)
username=danjac                        ; (default is no username (open server))
password=m0nk3Y                           ; (default is no password (open server))

[supervisord]
logfile=/home/danjac354/logs/user/supervisor/supervisord.log       ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB                             ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                                ; (num of main logfile rotation backups;default 10)
loglevel=debug                                    ; (log level;default info; others: debug,warn,trace)
pidfile=/home/danjac354/supervisord.pid                                ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                                    ; (start in foreground if true;default false)
minfds=1024                                       ; (min. avail startup file descriptors;default 1024)
minprocs=200                                      ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:25863
username=danjac
password=m0nk3Y

[program:nginx]
directory=/home/danjac354/webapps/nginx
command=sbin/nginx 
user=danjac354
autostart=true
autorestart=true
logfile=/home/danjac354/logs/user/supervisor/nginx.log
redirect_stderr=true

[program:rentalerts]
directory=/home/danjac354/webapps/rentalerts/rentalerts
command=/home/danjac354/.virtualenvs/rentalerts/bin/python2.7 /home/danjac354/webapps/rentalerts/rentalerts/bin/production run_gunicorn -b 127.0.0.1:53369 -w 4
user=danjac354
autostart=true
autorestart=true
logfile=/home/danjac354/logs/user/supervisor/rentalerts.log
redirect_stderr=true

